# apps/home-assistant/claude-stats/com.claude.stats-pusher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.stats-pusher</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source ~/.claude-stats.env &amp;&amp; /usr/local/bin/claude-stats-push</string>
    </array>
    
    <key>StartInterval</key>
    <integer>900</integer><!-- Every 15 minutes -->
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/claude-stats-push.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/claude-stats-push.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
